# uki-demo/actor-demo.ini
[program:actor-demo]
command = /root/venv36/bin/dramatiq --processes 1 task --queues actor-demo
directory = /root/work/uki-demo/
redirect_stderr = true
stdout_logfile = /data/logs/console/actor-demo.log
autostart = true
autorestart = true
startsecs = 10
stopwaitsecs = 600
environment = ENV='pro'