# uki-demo/uki-demo.ini
[program:uki-demo]
command = python3 main.py
directory = /root/work/uki-demo/
redirect_stderr = true
stdout_logfile = /data/logs/console/uki-demo.log
autostart = true
autorestart = true
startsecs = 10
stopwaitsecs = 600
environment = ENV='pro'